# supervisor/pelican.conf
[program:pelican]
command=python -m pelican.server
directory=/usr/local/blog/output
autostart=true
autorestart=true
